# doc/gunicorn.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='application/pyconsk' type='service' version='0'>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <exec_method name='start' type='method' exec='%{daemon/bin} -c %{daemon/config} %{daemon/app}' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    <exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='60'/>
    <instance name='live' enabled='true'>
      <method_context working_directory='/data/www/pycon.sk'>
        <method_credential group='pyconsk' user='pyconsk'/>
        <method_environment>
          <envvar name='PYTHONPATH' value='/data/www/pycon.sk'/>
          <envvar name='VIRTUAL_ENV' value='/data/www/pycon.sk/envs3'/>
          <envvar name='PATH' value='/data/www/pycon.sk/envs3/bin:/usr/local/bin:/opt/local/bin:/usr/bin:/bin'/>
        </method_environment>
      </method_context>
      <property_group name='daemon' type='application'>
        <propval name='app' type='astring' value='pycon_sk.wsgi'/>
        <propval name='bin' type='astring' value='/data/www/pycon.sk/envs3/bin/gunicorn'/>
        <propval name='config' type='astring' value='/data/www/pycon.sk/doc/gunicorn.py'/>
      </property_group>
    </instance>
    <instance name='staging' enabled='true'>
      <method_context working_directory='/data/www/staging.pycon.sk'>
        <method_credential group='pyconsk' user='pyconsk'/>
        <method_environment>
          <envvar name='PYTHONPATH' value='/data/www/staging.pycon.sk'/>
          <envvar name='VIRTUAL_ENV' value='/data/www/staging.pycon.sk/envs3'/>
          <envvar name='PATH' value='/data/www/staging.pycon.sk/envs3/bin:/usr/local/bin:/opt/local/bin:/usr/bin:/bin'/>
        </method_environment>
      </method_context>
      <property_group name='daemon' type='application'>
        <propval name='app' type='astring' value='pycon_sk.wsgi'/>
        <propval name='bin' type='astring' value='/data/www/staging.pycon.sk/envs3/bin/gunicorn'/>
        <propval name='config' type='astring' value='/data/www/staging.pycon.sk/doc/staging.py'/>
      </property_group>
    </instance>
    <stability value='Evolving'/>
  </service>
</service_bundle>
